# scripts/linstor-controller.service
[Unit]
Description=LINSTOR Controller Service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/linstor-server
ExecStart=/opt/linstor-server/bin/Controller
User=root

[Install]
Wanted=multi-user.target
